# syslog/supervisor/syslog-beward-ds.conf
[program:syslog-beward-ds]
process_name = %(program_name)s_%(process_num)02d
command = node /opt/rbt/syslog/beward.js --config=beward_ds
autostart = true
autorestart = true
numprocs = 1
redirect_stderr = true
stdout_logfile = /var/log/syslog_intercom/beward-ds.log
stderr_logfile = /var/log/syslog_intercom/beward-ds-error.log
stopwaitsecs = 3600
startretries = 5
