# server/docker/supervisor/default-task.ini
[program:default-task]
process_name=%(program_name)s_%(process_num)02d
command=php /srv/app/task.php --queue=default
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/srv/app/var/log/default.log
stderr_logfile=/srv/app/var/log/default-error.log
stopwaitsecs=3600
startretries=5
